# tlp-dyn-saver.service.in
# TLP Dynamic SAV Daemon
# Adjusts power-saver profile parameters based on workload lag detection
# Managed by TLP (tlp.in) - NOT automatically started
#
# Copyright (c) 2026 Mario Herrmann <mameogdoay at posteo.de> and others.
# SPDX-License-Identifier: GPL-2.0-or-later

[Unit]
Description=TLP Dynamic SAV Daemon
After=tlp.service
PartOf=tlp.service

[Service]
Type=simple
ExecStart=@TLP_TLIB@/func.d/35-tlp-func-dyn-saver run
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal

# Lockdown
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/tlp /tmp
ProtectClock=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectKernelLogs=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
KeyringMode=private

